# profile.d/configure.sh
#!/bin/bash
HOMEDIR=$HOME
sudo ln -s $HOMEDIR/.dotfiles/profile.d/shared_default_variables.sh /etc/profile.d/
