# profile.d/shared_default_variables.sh
export EDITOR=vim
export TERMINAL=tilix
